# deps/redns.service
[Unit]
Description=redns DNS IOC

[Service]
ExecStart=/usr/bin/rkt run \
	--port=dns-tcp:8053 \
	--port=dns-udp:8053 \
	--volume redis,kind=host,source=/tmp/redns_redis/ \
	unixvoid.com/redis --mount volume=redis,target=/redisbak/ \
        unixvoid.com/redns

#[Service]
#ExecStart=/usr/bin/rkt run \
#	--net=host \
#	--volume redis,kind=host,source=/tmp/redns_redis/ \
#	unixvoid.com/redis --mount volume=redis,target=/redisbak/ \
#       unixvoid.com/redns

[Install]
WantedBy=multi-user.target
